# deploy/systemd/blog-backend.service
[Unit]
Description=Blog Platform Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/blog/backend

# 环境变量
Environment="GIN_MODE=release"
EnvironmentFile=-/var/www/blog/backend/.env

# 执行命令
ExecStart=/var/www/blog/backend/bin/server

# 重启策略
Restart=always
RestartSec=10s

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blog-backend

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
